# templates/glpi.cron.d
*/1 * * * * www-data /usr/bin/php5 /usr/share/glpi/front/cron.php
